# root/etc/s6-overlay/s6-rc.d/pre-qbittorrent/up
s6-setuidgid hotio /etc/s6-overlay/s6-rc.d/pre-qbittorrent/configure-qbittorrent-pre.sh
