# root/etc/s6-overlay/s6-rc.d/post-qbittorrent/up
/etc/s6-overlay/s6-rc.d/post-qbittorrent/configure-qbittorrent.sh
